# build_files/services/games.service
[Unit]
Description="Decrypts and mount games disks"

[Service]
ExecStart=sh /usr/lib/systemd/system/games.sh
Type=oneshot
RemainAfterExit=yes
