# build_files/services/games.timer
[Unit]
Description="Timer for Games disk decryption and mounting"

[Timer]
OnBootSec=5s
Unit=games.service

[Install]
WantedBy=multi-user.target
